# applications/all-in-one/supevirsord/composite.conf
[inet_http_server]
port=127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001

[supervisord]

[program:all-in-one]
environment = PYTHONUNBUFFERED=1
directory=/var/supervisord/all-in-one
command=gunicorn -w 4 -k uvicorn.workers.UvicornWorker all_in_one.server:app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisord/all-in-one.err.log
stdout_logfile=/var/log/supervisord/all-in-one.out.log
